# N1/files/etc/crontabs/root
0 4 * * * fstrim -a
0 4 * * * /mnt/sda/aria2/tracker.sh
